# roles/debian_base/files/etc/systemd/system/bind9.service.d/override.conf
[Unit]
BindsTo=sys-subsystem-net-devices-ffhl.device
After=network.target sys-subsystem-net-devices-ffhl.device
