# roles/meshvpn/files/systemd/update-ffhl-mesh-vpn.service
[Service]
Type=oneshot
#WorkingDirectory=/var/local/ffhl-mesh-vpn-peers
ExecStart=/usr/local/lib/ffhl/update-meshvpn-keys.sh
