# roles/base/files/etc/systemd/system/fastd@.service.d/override.conf
[Unit]
PartOf=fastd.service

[Service]
# make sure these dirs exists for fastd dignostics/metrics socket
ExecStartPre=mkdir -p /run/fastd
ExecStartPre=chown fastd:fastd /run/fastd
