# roles/debian_base/files/etc/systemd/system/fastd@.service.d/override.conf
[Service]
ExecStartPre=-/usr/bin/mkdir /run/fastd
ExecStartPre=/usr/bin/chown fastd:fastd /run/fastd
