# roles/base/files/etc/systemd/system/update-icvpn.service
[Service]
Type=oneshot
WorkingDirectory=/etc/tinc/icvpn
ExecStart=/usr/bin/git pull
