# roles/base/files/etc/systemd/system/update-ffhl-dns.service
[Service]
Type=oneshot
WorkingDirectory=/var/local/ffhl-dns
ExecStart=/usr/bin/git pull
