# roles/icvpn/files/etc/systemd/system/update-icvpn.timer
[Timer]
OnCalendar=*:00/16

[Install]
WantedBy=multi-user.target
