# roles/icvpn/files/etc/systemd/system/update-icvpn-meta.service
[Service]
Type=oneshot
WorkingDirectory=/var/lib/icvpn-meta
ExecStart=/usr/bin/git pull
