# roles/ffhl_nameserver/files/systemd/update-ffhl-dns.service
[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/local/lib/ffhl/update-dns.sh
