# roles/ffhl_ns_auth/files/systemd/update-ffhl-dns.service
[Service]
Type=oneshot
ExecStart=/usr/local/lib/ffhl/update-dns.sh