# roles/firmware_mirror/files/etc/systemd/system/firmware-mirror.timer
[Unit]
Description=Firmware Mirror Timer

[Timer]
OnCalendar= *-*-* 03:42:23

[Install]
WantedBy=timers.target