# roles/common_gateway/files/etc/cron.d/update-apt-list
# update the package lists so prometheus can alert us
# if there are many updates available

7 22 * * * apt-get update
