# roles/old-arch/files/etc/systemd/system/update-ffhl-mesh-vpn.service
[Service]
Type=oneshot
WorkingDirectory=/etc/fastd/ffhl-mesh-vpn/peers
ExecStart=/usr/bin/git pull
