# roles/old-arch/files/etc/systemd/system/update-ffhl-dns.service
[Service]
Type=oneshot
WorkingDirectory=/var/named/ffhl
ExecStart=/usr/bin/git pull
