# roles/base/files/etc/udev/rules.d/90-fix-gre.rules
ACTION=="add", KERNEL=="ffrhein-*", SUBSYSTEM=="net", RUN+="/usr/bin/ip link set %k mtu 1400"
