# roles/base/files/mesh-vpn/update-ffhl-mesh-vpn.service
[Service]
Type=oneshot
WorkingDirectory=/var/local/ffhl-mesh-vpn-peers
ExecStart=/usr/bin/git pull
