# roles/debian_base/files/etc/systemd/system/update-ffhl-mesh-vpn.service
[Service]
Type=oneshot
WorkingDirectory=/etc/fastd/ffhl_mesh_vpn/peers
ExecStart=/usr/bin/git pull
