# roles/debian_base/files/etc/systemd/system/update-ffhl-dns.service
[Service]
Type=oneshot
WorkingDirectory=/var/cache/bind
ExecStart=/usr/bin/git pull
